# profile.d/00-lua-lsp.rc
export PATH="$PATH:$HOME/lua-language-server/bin"
